# crontab
* * * * * bash /root/prepare-minitaf-server/weekly-maintenance.sh
